# crontab/dockeruser
@reboot curl http://bmc:5000/bmc/update/${WEBHOOK}
20 10 * * * curl http://bmc:5000/bmc/update/${WEBHOOK}
